# raspberry-pi/bluezcript-webui.service
[Unit]
Description=BlueZcript Web UI
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/Bluezcript/raspberry-pi
ExecStart=/usr/bin/python3 /home/pi/Bluezcript/raspberry-pi/web_app.py
Restart=always
RestartSec=10

# Environment
Environment="PATH=/usr/bin:/usr/local/bin"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bluezcript-webui

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
